# deploy/systemd/defensive-ai-gateway.service
[Unit]
Description=Defensive AI Gateway
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/defensive-ai-gateway
ExecStart=/usr/bin/python3 -m defensive_ai_gateway --config /etc/defensive-ai-gateway/prod.yaml
Restart=always
RestartSec=3
User=defensive-ai
Group=defensive-ai
Environment=DEFENSIVE_AI_DB=/var/lib/defensive-ai-gateway/gateway.db
EnvironmentFile=/etc/defensive-ai-gateway/env
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
TasksMax=256
MemoryMax=2G
ReadWritePaths=/var/lib/defensive-ai-gateway

[Install]
WantedBy=multi-user.target
